# Configs/Public/systemd/user/syncthing.service
[Unit]
Description=Syncthing
After=network.target

[Service]
ExecStart=/usr/bin/syncthing
Restart=always

[Install]
WantedBy=default.target
